# ansible/playbooks/templates/rc.local.j2
#!/bin/bash

ifconfig eth1 {{ eth1_ip }}/24 up
